# postgres/crontab
0 0 * * * /var/lib/postgresql/backup.sh >> /var/log/backup.log
